# .crontab
# crontab
# +--------- minute (0-59)                                              | output dumper: >/dev/null 2>&1
# |     +------- hour (0-23)                                            | multiple values use commas: 3,12,47
# |     |       +----- day of month (1-31)                              | do every x intervals: */x
# |     |       |       +--- month (1-12)                               | aliases: @reboot @hourly
# |     |       |       |       +- day of week (0-6) (sunday = 0)       |          @weekly @monthly
# |     |       |       |       |                                       |          @yearly
# *     *       *       *       *       command                         |

*/15    *       *       *       *       /home/jam/bin/wal-set

# vim:set ft=crontab
